# projects/Amlogic-ce/packages/sysutils/amremote/system.d/amlogic-remote-toggle.service
[Unit]
Description=Remote toggle service for Amlogic devices
Before=graphical.target kodi.service
After=multi-user.target
ConditionPathExists=/proc/device-tree/meson-ir/compatible
ConditionPathExists=/proc/device-tree/meson-remote/compatible
ConditionPathIsDirectory=!/proc/device-tree/custom_maps

[Service]
Type=oneshot
ExecStart=/bin/sh -c "/usr/lib/coreelec/remote-toggle reboot"

[Install]
WantedBy=graphical.target
