# packages/sx05re/emuelec/system.d/emuelec-disable_small_cores.service
[Unit]
Description=emuelec disable small cores service
Before=emulationstation.service
After=graphical.target

[Service]
Type=oneshot
ExecStart=/usr/bin/emuelec-utils small-cores disable
RemainAfterExit=true
ExecStop=/usr/bin/emuelec-utils small-cores enable

[Install]
WantedBy=emuelec.target
