# packages/addons/service/tvheadend44/source/system.d/service.tvheadend44.service
[Unit]
Description=TVHeadend44 Service
After=network-online.service
Requires=network-online.service

[Service]
ExecStart=/bin/sh -c "exec sh /storage/.kodi/addons/service.tvheadend44/bin/tvheadend44.start"
TimeoutStopSec=2
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=kodi.target
