# projects/Amlogic-ce/packages/sysutils/tee_preload_fw/system.d/video-firmware-preload.service
[Unit]
Description=Amlogic OPTEE Video Firmware Preload Service
After=opentee_linuxdriver.service
Requires=opentee_linuxdriver.service
Before=kodi.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecCondition=/bin/sh -c 'systemctl is-active --quiet opentee_linuxdriver.service && exit 0 || exit 1'
ExecStartPre=/usr/sbin/trusted-application-setup
ExecStart=/usr/sbin/tee_preload_fw /lib/firmware/video/video_ucode.bin

[Install]
WantedBy=multi-user.target
