# packages/sx05re/tools/sysutils/sixaxis/udev.d/99-sixaxis.rules
ACTION=="add", SUBSYSTEMS=="input", ATTRS{name}=="*PLAYSTATION(R)3 Controller", TAG+="systemd", ENV{SYSTEMD_WANTS}="sixaxis@%E{DEVNAME}.service", ENV{ID_INPUT_KEY}=""
ACTION=="add", SUBSYSTEMS=="input", ATTRS{name}=="*PLAYSTATION(R)3Conteroller", TAG+="systemd", ENV{SYSTEMD_WANTS}="sixaxis@%E{DEVNAME}.service", ENV{ID_INPUT_KEY}=""
ACTION=="add", SUBSYSTEMS=="input", ATTRS{name}=="*PLAYSTATION(R)3Conteroller-PANHAI", TAG+="systemd", ENV{SYSTEMD_WANTS}="sixaxis@%E{DEVNAME}.service", ENV{ID_INPUT_KEY}=""
ACTION=="add", SUBSYSTEMS=="input", ATTRS{name}=="*PS(R) Gamepad", TAG+="systemd", ENV{SYSTEMD_WANTS}="sixaxis@%E{DEVNAME}.service", ENV{ID_INPUT_KEY}=""
ACTION=="add", SUBSYSTEMS=="input", ATTRS{name}=="*PS3 GamePad", TAG+="systemd", ENV{SYSTEMD_WANTS}="sixaxis@%E{DEVNAME}.service", ENV{ID_INPUT_KEY}=""
ACTION=="add", SUBSYSTEMS=="input", ATTRS{name}=="Sony Interactive Entertainment Wireless Controller", TAG+="systemd", ENV{SYSTEMD_WANTS}="sixaxis@%E{DEVNAME}.service", ENV{ID_INPUT_KEY}=""
ACTION=="add", SUBSYSTEMS=="input", ENV{DEVNAME}=="/dev/input/event*", ATTRS{name}=="*Motion Sensors", RUN+="/bin/rm %E{DEVNAME}", ENV{ID_INPUT_JOYSTICK}="", ENV{ID_INPUT_KEY}=""
