# projects/Amlogic-ce/packages/sysutils/dtb-xml/system.d/dtb-xml.service
[Unit]
Description=Copy user edid to dtb

[Service]
Type=oneshot
RemainAfterExit=true
ExecStop=-/usr/lib/coreelec/dtb-xml -e

[Install]
WantedBy=multi-user.target
