# packages/sx05re/libretro/retroarch/system.d/tmp-joypads.mount
[Unit]
Description=Joypad configs directory
After=systemd-tmpfiles-setup.service
DefaultDependencies=no

[Mount]
What=none
Where=/tmp/joypads
Type=overlay
Options=lowerdir=/etc/retroarch-joypad-autoconfig,upperdir=/storage/joypads,workdir=/storage/.tmp/joypads-workdir

[Install]
WantedBy=emuelec.target
