# packages/third/tools/acelist/system.d/acelist.timer
[Timer]
OnBootSec=2min
OnUnitActiveSec=50min
Unit=acelist.service
[Install]
WantedBy=kodi.target
