# packages/lib32/tools/lib32-ldconfig/system.d/ldconfig-watch.path
# SPDX-License-Identifier: GPL-2.0-or-later
# Copyright (C) 2022-present 7Ji (https://github.com/7Ji)

[Unit]
Description=Rebuild Dynamic Linker Cache If EmuELEC Libs Updated

[Path]
PathChanged=/emuelec/lib32
PathChanged=/emuelec/lib

[Install]
WantedBy=multi-user.target