# projects/Amlogic-ce/devices/Amlogic-ne/packages/sysutils/inject_bl301/system.d/update-bl301.service
[Unit]
Description=CoreELEC BL301 Update Service

[Service]
Type=oneshot
ExecStart=/usr/lib/coreelec/update-bl301 -v
RemainAfterExit=no

[Install]
WantedBy=basic.target
