# projects/Amlogic-ce/devices/Amlogic-no/packages/mediacenter/CoreELEC-settings/system.d/backup-restore.service
[Unit]
Description=Restoring Backup
DefaultDependencies=no

[Service]
Type=idle
ExecStart=/usr/lib/coreelec/backup-restore
StandardInput=tty-force
StandardOutput=inherit
StandardError=inherit
