# packages/linux-firmware/amlogic/rtkbt-firmware-aml/system.d/rtkbt-firmware-aml.service
[Unit]
Description=Attach /dev/ttyS1 to BlueZ stack
Wants=bluetooth.service
Before=bluetooth.service
After=dev-ttyS1.device

[Service]
Type=simple
ExecStartPre=/bin/sh -c 'echo 1 > /sys/class/rfkill/rfkill0/state'
ExecStartPre=/usr/sbin/rfkill unblock bluetooth
ExecStart=/usr/bin/rtk_hciattach -n -s 115200 /dev/ttyS1 rtk_h5
ExecStopPost=/usr/sbin/rfkill block bluetooth
