# packages/audio/brutefir/config/system.d/service.brutefir.service.sample
[Unit]
Description=Brutefir
Wants=kodi.service
After=kodi.service

[Service]
ExecStart=/usr/bin/brutefir /storage/.config/.brutefir_config
Restart=always

[Install]
WantedBy=kodi.service
