# packages/sysutils/amremote/system.d/amlogic-remote-toggle.service
[Unit]
Description=Remote toggle service for Amlogic devices
After=local-fs.target
ConditionPathExists=|/proc/device-tree/mali@d00c0000/compatible
ConditionPathExists=|/proc/device-tree/t82x@d00c0000/compatible

[Service]
Type=oneshot
ExecStart=/bin/sh -c "/usr/lib/coreelec/remote-toggle reboot"

[Install]
WantedBy=multi-user.target
