# packages/sysutils/amremote/system.d/amlogic-remotecfg.service
[Unit]
Description=Amlogic IR remote support
ConditionPathExists=/dev/amremote

[Service]
Type=oneshot
ExecStart=/usr/lib/coreelec/remote-config
RemainAfterExit=yes

[Install]
WantedBy=basic.target
