# packages/sx05re/tools/sound/fluidsynth-git/system.d/fluidsynth.service
[Unit]
Description=FluidSynth Software Synthesizer Daemon
After=sound.target

[Service]
EnvironmentFile=/etc/fluidsynth/fluidsynth.conf
ExecStart=/usr/bin/fluidsynth -is -a $AUDIO_DRIVER -m $MIDI_DRIVER -p $PORT_NAME $OTHER_OPTS $SOUND_FONT

[Install]
WantedBy=multi-user.target
