# projects/Rockchip/devices/OdroidGoAdvance/packages/odroidgoa-utils/system.d/odroidgoa-hotkeys.service
[Unit]
Description=Odroid Go Advance Hotkeys
ConditionPathExists=/dev/input/event0
ConditionPathExists=/dev/input/event2
Before=emuelec.target

[Service]
Environment=HOME=/storage
EnvironmentFile=/storage/.config/emulationstation/scripts/es_env.sh
ExecStart=/usr/bin/oga_events.py
Type=simple

[Install]
WantedBy=multi-user.target