# projects/Amlogic/filesystem/usr/lib/systemd/system/kodi-aml-hdmimonitor.service
[Unit]
Description=Kodi Amlogic HDMI connection monitor
After=kodi.service

[Service]
Type=simple
ExecStart=/usr/lib/kodi/aml-hdmimonitor.sh

[Install]
WantedBy=kodi.target
