# packages/linux-drivers/amlogic/qca9377-aml/udev.d/99-qca9377-wireless.rules
ACTION=="add|change", SUBSYSTEM=="net", ENV{INTERFACE}=="wlan*", \
  ATTRS{vendor}=="0x0271", ATTRS{device}=="0x070?", \
  RUN+="/usr/bin/systemctl restart connman.service"
