# packages/network/nfs-utils/system.d/load-module-nfsd.service
[Unit]
Description=Load module nfsd service

[Service]
Type=oneshot
ExecStart=/usr/sbin/modprobe nfsd
RemainAfterExit=true
