# projects/Amlogic-ce/packages/linux-drivers/openvfd-driver/system.d/openvfd.service
[Unit]
Description=OpenVFD Service
ConditionPathExists=/proc/device-tree/openvfd/
After=local-fs.target

[Service]
ExecStart=/usr/lib/coreelec/openvfd-start
ExecStop=/bin/kill -TERM $MAINPID
ExecStopPost=-/usr/sbin/rmmod openvfd
RemainAfterExit=yes

[Install]
WantedBy=basic.target
