# projects/Amlogic-ce/devices/Amlogic-ne/packages/linux-drivers/amlogic/media_modules-aml/system.d/media_modules-aml.service
[Unit]
Description=Amlogic media_modules-aml service
Before=kodi.service
After=opentee_linuxdriver.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/coreelec/media_modules-aml

[Install]
WantedBy=multi-user.target
