# projects/Amlogic-ce/packages/linux-drivers/amlogic/opentee_linuxdriver/system.d/opentee_linuxdriver.service
[Unit]
Description=Optee TEE-Supplicant and Video Firmware Preload service
Before=kodi.service
After=local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/coreelec/tee-loader start
ExecStart=/usr/lib/coreelec/dovi-loader start
ExecStop=/usr/lib/coreelec/dovi-loader stop
ExecStop=/usr/lib/coreelec/tee-loader stop

[Install]
WantedBy=multi-user.target
