# packages/sx05re/tools/sixaxis/system.d/sixaxis@.service
[Unit]
Description=sixaxis helper (%I)

[Service]
Type=simple
ExecStart=/usr/bin/sixaxis-helper.sh %I

[Install]
WantedBy=multi-user.target
