# projects/Amlogic-ce/packages/linux-firmware/amlogic/qca-firmware-aml/system.d/qca-firmware-aml.service
[Unit]
Description=Attach /dev/ttyS1 to BlueZ stack using qca type
Wants=bluetooth.service
Before=bluetooth.service
After=dev-ttyS1.device
StartLimitInterval=60
StartLimitBurst=2

[Service]
Type=simple
ExecStartPre=/usr/sbin/rfkill unblock bluetooth
ExecStart=/usr/bin/hciattach -n -s 115200 /dev/ttyS1 qca 2000000
ExecStopPost=/usr/sbin/rfkill block bluetooth
Restart=always
RestartSec=2
LimitNPROC=1
TimeoutStopSec=1s
